# root/etc/hotplug.d/usb/20-usbleach
# @TODO: send data to log = call module via a cli interface
